# backend/scripts/stemscribe-queue-monitor.timer
[Unit]
Description=Run StemScriber queue-depth monitor every 5 minutes
Requires=stemscribe-queue-monitor.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=stemscribe-queue-monitor.service

[Install]
WantedBy=timers.target
